# crontab
#!/bin/sh
curl -i localhost:5000/notify/global-repo-dimension
